# CreditUpdate.cron
00 10 * * * /home/ubuntu/kuzo/venv/bin/python3 /home/ubuntu/kuzo/manage.py subscription >> /home/ubuntu/kuzo/cron.log 2>&1
