# Script.cron
2 * * * * source venv/bin/activate && python manage.py transfer_funds
